# ad1000.conf
#!upstart
description "ad1000"

env USER=root
env DAEMON="/usr/local/bin/ad1000"
env PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

start on started xbmc 
stop on runlevel [06]

exec ionice -c best-effort $DAEMON > /dev/null || :
